# motion/crontab
* * * * * git/rpi/motion/tmux-motion.sh

