# weather2/crontab
* * * * * /home/pi/git/rpi/weather2/weather.py 12761582 2 1 /home/pi/weather.txt 600

