# weather/crontab
* * * * * /home/pi/git/rpi/weather/weather.py 12761450 2 1 /home/pi/weather.txt 600

