# weather/crontab
* * * * * /home/pi/git/rpi/weather/weather.py 12761582 2 1 /home/pi/weather.txt 300

