# clock/crontab
* * * * * /home/pi/git/rpi/clock/clock.py 12761450 2 1 /home/pi/clock.txt 600

