# net/rc.local
#!/bin/sh -e

/usr/local/bin/pigpiod

exit 0
